# systemd/hermes-flight-recorder-capture.service.in
[Unit]
Description=Capture Hermes events into the Flight Recorder outbox

[Service]
Type=oneshot
UMask=0077
# Bound a hung pass so it goes 'failed' and the next tick re-arms, instead of
# blocking forever (the default is infinity).
TimeoutStartSec=120
ExecStart=@RT@ run --flight-recorder-home @FR_HOME@ --hermes-home @HERMES_HOME@
