# systemd/hermes-flight-recorder-reconcile.timer.in
[Unit]
Description=Reconcile Hermes events every minute

[Timer]
# Realtime anchor, once a minute (at :00). deploy.sh renders this unit but
# preserves its enabled/disabled state — reconcile ships disabled by default.
OnCalendar=*:*:00
Persistent=true
AccuracySec=5s
Unit=hermes-flight-recorder-reconcile.service

[Install]
WantedBy=timers.target
